# pytemp.conf
## pytemp

description "PythonTemps"
author "J Archer"
 
start on net-device-up IFACE!=lo
 
respawn
 
exec /usr/bin/python /usr/local/bin/pytemp.py

